# rc.d/g_guest_segments
#!/bin/sh
#
# PROVIDE: g-guest_segments
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

# IMPORTANT: pas de '-' dans 'name'
name="g_guest_segments"
rcvar="g_guest_segments_enable"

: ${g_guest_segments_enable:="NO"}
: ${g_guest_segments_base:="/usr/local/g-guest"}
: ${g_guest_segments_cmd:="${g_guest_segments_base}/bin/segments.sh"}

start_cmd="g_guest_segments_start"
stop_cmd="g_guest_segments_stop"
status_cmd="g_guest_segments_status"
extra_commands="plan"
plan_cmd="g_guest_segments_plan"

g_guest_segments_start()
{
    if [ ! -x "${g_guest_segments_cmd}" ]; then
        echo "g-guest_segments: ${g_guest_segments_cmd} introuvable ou non exécutable" >&2
        return 1
    fi
    echo "g-guest_segments: ensure (lecture seule)"
    BASE_DIR="${g_guest_segments_base}" \
    CFG_DIR="${g_guest_segments_base}/config" \
    "${g_guest_segments_cmd}" ensure
}

g_guest_segments_stop()
{
    echo "g-guest_segments: rien à arrêter (one-shot)."
    return 0
}

g_guest_segments_status()
{
    if [ ! -x "${g_guest_segments_cmd}" ]; then
        echo "g-guest_segments: ${g_guest_segments_cmd} introuvable ou non exécutable" >&2
        return 1
    fi
    echo "g-guest_segments: état (ensure, lecture seule)"
    BASE_DIR="${g_guest_segments_base}" \
    CFG_DIR="${g_guest_segments_base}/config" \
    "${g_guest_segments_cmd}" ensure
}

g_guest_segments_plan()
{
    if [ ! -x "${g_guest_segments_cmd}" ]; then
        echo "g-guest_segments: ${g_guest_segments_cmd} introuvable ou non exécutable" >&2
        return 1
    fi
    echo "g-guest_segments: plan (texte)"
    BASE_DIR="${g_guest_segments_base}" \
    CFG_DIR="${g_guest_segments_base}/config" \
    "${g_guest_segments_cmd}" plan
}

load_rc_config "$name"
run_rc_command "$1"
